# config/files/usr/etc/xdg/autostart/naia-shell.desktop
[Desktop Entry]
Name=Naia
Comment=Personal AI avatar shell
Exec=sh -c 'flatpak info io.nextain.naia >/dev/null 2>&1 && flatpak run io.nextain.naia || true'
Icon=io.nextain.naia
Type=Application
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
