# config/files/usr/etc/xdg/autostart/naia-shell.desktop
[Desktop Entry]
Name=Naia
Comment=Personal AI avatar shell
Exec=flatpak run io.nextain.naia
Icon=naia
Type=Application
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
