# config/files/usr/lib/systemd/user/naia-gateway.service
[Unit]
Description=Naia OpenClaw Gateway
Documentation=https://docs.openclaw.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/naia-gateway-wrapper
Restart=on-failure
RestartSec=5
Environment=OPENCLAW_PROFILE=naia
Environment=NODE_ENV=production

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=%h/.naia
PrivateTmp=true

[Install]
WantedBy=default.target
