# config/systemd/naia-gateway.service
[Unit]
Description=Naia Gateway (OpenClaw)
Documentation=https://github.com/nichochar/openclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/.naia/openclaw/node_modules/openclaw/openclaw.mjs gateway run --bind loopback --port 18789 --allow-unconfigured
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.config/nvm/versions/node/v22/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=%h/.naia/openclaw
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.naia %h/.openclaw
PrivateTmp=true

[Install]
WantedBy=default.target
