# config/systemd/naia-gateway.service
[Unit]
Description=Naia Gateway (optional external OpenClaw gateway)
Documentation=https://github.com/nextain/naia-os
After=network-online.target
Wants=network-online.target

# NOTE: Naia Agent (TypeScript) is the built-in runtime — this service is
# optional and only needed if you run an external OpenClaw-compatible gateway.
# Configure the gateway URL in Naia Shell → Settings → Gateway URL.
#
# To use this service, install openclaw separately and update ExecStart below:
# ExecStart=/usr/bin/node /path/to/openclaw/openclaw.mjs gateway run --bind loopback --port 18789 --allow-unconfigured

[Service]
Type=simple
# ExecStart is intentionally left unconfigured — Naia Agent runs inside Naia Shell.
# Uncomment and adjust the line below if you use an external OpenClaw gateway:
# ExecStart=/usr/bin/node %h/.naia/openclaw/node_modules/openclaw/openclaw.mjs gateway run --bind loopback --port 18789 --allow-unconfigured
ExecStart=/bin/true
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.config/nvm/versions/node/v22/bin:/usr/local/bin:/usr/bin:/bin
Restart=no

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.naia %h/.openclaw
PrivateTmp=true

[Install]
WantedBy=default.target
